# mast3r-camera.service
[Unit]
Description=MASt3R Camera Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/mast3r-camera-client
ExecStart=/usr/bin/python3 /home/pi/mast3r-camera-client/camera_client.py --host linux-2 --port 5050 --fps 1
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
